# SaveGuardian/SaveGuardian.service
[Unit]
Description=Save Guardian

[Service]
Type=notify
ExecStart=/usr/sbin/SaveGuardian

[Install]
WantedBy=multi-user.target